# data/units/cldd_rw.socket.in
[Socket]
ListenStream=@RUNDIR@/cldd/socket
SocketMode=0600
Service=cldd.service

[Install]
WantedBy=sockets.target
