# data/units/cldd_ro.socket.in
[Socket]
ListenStream=@RUNDIR@/cldd/socket_ro
SocketMode=0400
Service=cldd.service

[Install]
WantedBy=sockets.target
